# config/supervisor/crontab.ini
[program:crontab]
command=/usr/sbin/crond start
autostart=true